# lila-tablebase.socket
[Unit]
Description=Tablebase server socket

[Socket]
ListenStream=127.0.0.1:9000

[Install]
WantedBy=sockets.target
