# lila-tablebase.socket
[Unit]
Description=Tablebase server socket

[Socket]
ListenStream=/run/lila-tablebase.socket

[Install]
WantedBy=sockets.target
